# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx-access.log
stderr_logfile=/var/log/supervisor/nginx-error.log

[program:uvicorn]
command=uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 4
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/uvicorn-access.log
stderr_logfile=/var/log/supervisor/uvicorn-error.log
environment=PATH="/usr/local/bin:%(ENV_PATH)s"
